# claudebot/com.emprezario.claudebot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.emprezario.claudebot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/edduran/Documents/GitHub/ai-realtor/claudebot/main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/edduran/Documents/GitHub/ai-realtor/claudebot</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/edduran/.local/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/edduran/.claudebot/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/edduran/.claudebot/logs/stderr.log</string>
</dict>
</plist>
